# supervisor/script.conf
[supervisord]
nodaemon=true
user=root

[program:docker_demon]
command= dockerd
autostart=true
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0

[program:rod_manager]
command= sudo docker run -p 7317:7317 ghcr.io/go-rod/rod
autostart=true
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0


[program:gosandbox_remote]
command=  reflex -g '*.go' go run main.go --start-service
autostart=true
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0